# etc/profile.d/aux-cli-completion.sh
#!/bin/bash
# _auxcli_completions()
# {
#     . /bin/aux-cli help  >/dev/null 2>&1
#     COMPREPLY=($(compgen -W "${commands}" "${COMP_WORDS[1]}"))
# }

# complete -F _auxcli_completions aux-cli